# packaging/nosignal-daemon.service
# NoSignal daemon — systemd user unit.
# Installed by the .deb/.rpm to /usr/lib/systemd/user/.
# Enable with: systemctl --user enable --now nosignal-daemon
#

[Unit]
Description=NoSignal display daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=io.github.shadowxcc.NoSignal.Daemon1
ExecStart=/usr/bin/nosignald
Restart=on-failure
Slice=background.slice

[Install]
WantedBy=graphical-session.target
